# crontab
* * * * * root echo "Hello world" >> /var/log/cron.log
